# config/web-socket.conf
[program:websocket]
command=php /vagrant/bin/ws.php
stderr_logfile=/vagrant/var/logs/web-socket-stderr.log
stdout_logfile=/vagrant/var/logs/web-socket-stdout.log
